# gerrit-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:gerrit]
command=/bin/bash -c "/home/gerrit/conf-and-run-gerrit.sh"
